# image/etc/supervisord.d/rsyslogd.ini
[program:rsyslogd]
command=/sbin/rsyslogd -n
