# image/etc/supervisord.d/suricata.ini
[program:suricata]
command=/usr/sbin/suricata %(ENV_SURICATA_ARGS)s
stdout_logfile=syslog
stderr_logfile=syslog
