# image/etc/supervisord.d/syslog-tail.ini
[program:syslog-tail]
command=tail -f -n0 /var/log/messages
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
