# image/etc/supervisord.d/syslog-tail.ini
[program:syslog-tail]
command=tail -f /var/log/messages
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
